# oci_deploy/vfl-app.service
[Unit]
Description=VFL Prediction Engine — Flask app
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/vfl-deploy
EnvironmentFile=/etc/systemd/system/vfl-app.env
ExecStart=/opt/vfl-deploy/venv/bin/python railway_runner.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
